# systemd/user/bt-auto-connect.service
[Unit]
Description=auto connect bluetooth
BindsTo=niri.service
After=niri.service

[Service]
ExecStart=/home/jappe/.local/bin/bt-auto-connect

