# systemd/user/idle-inhibitor.service
[Unit]
Description=idle-inhibitor
BindsTo=niri.service
After=niri.service

[Service]
ExecStart=/home/jappe/.local/bin/idle-inhibitor
