# systemd/user/mako.service
[Unit]
Description=mako
BindsTo=niri.service
After=niri.service

[Service]
ExecStart=/usr/bin/mako
