# systemd/user/openrgb.service
[Unit]
Description=OpenRGB

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/openrgb --server --profile "1s-white-breathing"