# systemd/user/swww-overview.service
[Unit]
Description=swww daemon for the niri overview
BindsTo=niri.service
After=niri.service

[Service]
ExecStart=/usr/bin/swww-daemon --namespace niri-overview

