# systemd/user/waybar.service
[Unit]
Description=waybar
BindsTo=niri.service
After=niri.service

[Service]
ExecStart=/usr/bin/waybar
